# deploy/mcp-server.service
[Unit]
Description=Book Graph RAG MCP Server
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=gonzalo
Group=gonzalo
WorkingDirectory=/home/gonzalo/Gonzalo_codigo/Mcp_libro/MCP_neo4j_orangpi
ExecStart=/home/gonzalo/.local/bin/uv run book-graph-rag-mcp serve
Environment=PATH=/home/gonzalo/.local/bin:/usr/local/bin:/usr/bin:/bin
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
