# deploy/mcp-server.service
[Unit]
Description=Book Graph RAG MCP Server
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
# Replace 'bookgraph' with the actual user on the Orange Pi.
User=bookgraph
Group=bookgraph
WorkingDirectory=/home/bookgraph/Gonzalo_codigo/Mcp_libro/MCP_neo4j_orangpi
EnvironmentFile=/home/bookgraph/Gonzalo_codigo/Mcp_libro/MCP_neo4j_orangpi/.env
ExecStart=/home/bookgraph/.local/bin/uv run book-graph-rag-mcp serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
